# ramfox.service
[Unit]
Description=ramfox Service

[Service]
Type=simple
ExecStart=/usr/bin/ramfox

[Install]
WantedBy=multi-user.target